# templates/jellyfin.service.j2
[Unit]
Description=Jellyfin
After=network.target

[Service]
Type=simple
User={{ jellyfin_user }}
Restart=always
ExecStart=/bin/bash /usr/local/bin/jellyfin.sh

[Install]
WantedBy=multi-user.target