# templates/jellyfin.service.j2
[Unit]
Description=Jellyfin
After=network.target

[Service]
Type=simple
User={{ jellyfin_user }}
Restart=always
ExecStart={{ jellyfin_home }}/jellyfin.sh

[Install]
WantedBy=multi-user.target